# docker/php-fpm/crontab
* * * * * root echo "Hello world" >> /var/log/cron.log 2>&1
*\5 * * * * root php /app/bin/console app:update-courses > /var/log/cron.log
#