# docker/php-fpm/crontab
*/5 * * * * cd /app && /usr/local/bin/php bin/console app:update-courses > /var/log/command-log.txt -e prod
#